# deployment/task_scheduler_ocr_worker.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Windows Task Scheduler import file for PLW EDMS OCR Worker
  
  Import via:
    schtasks /Create /XML deployment\task_scheduler_ocr_worker.xml /TN "PLW-EDMS-OCR-Worker"
  
  Or: Task Scheduler GUI -> Action -> Import Task -> select this file
  
  Adjust:
    <UserId>     - Windows account running the service (use a dedicated service account)
    <WorkingDirectory> - Path to project root
    <Command>    - Full path to python.exe
    <Arguments>  - Full path to manage.py
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>PLW EDMS OCR Pipeline Worker - processes PDF queue every 30 seconds</Description>
    <Author>PLW-EDMS-Admin</Author>
  </RegistrationInfo>
  <Triggers>
    <BootTrigger>
      <Enabled>true</Enabled>
      <Delay>PT2M</Delay>
    </BootTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>PLW\edms_service</UserId>
      <LogonType>Password</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RestartOnFailure>
      <Interval>PT1M</Interval>
      <Count>10</Count>
    </RestartOnFailure>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Enabled>true</Enabled>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Python311\python.exe</Command>
      <Arguments>D:\plw_edms\manage.py run_ocr_worker --loop --interval 30</Arguments>
      <WorkingDirectory>D:\plw_edms</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
